# role-sync-bot.service.example
[Unit]
Description=Discord Role Sync Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/role-sync-bot
EnvironmentFile=/opt/role-sync-bot/.env
ExecStart=/opt/role-sync-bot/.venv/bin/python /opt/role-sync-bot/bot.py
Restart=on-failure
RestartSec=5
User=www-data

[Install]
WantedBy=multi-user.target
